# templates/etc/profile.d/postgres.sh.j2
export PGDATA={{ _postgresql_data_dir }}
export PATH=$PATH:{{ _postgresql_bin_path }}
